# scripts/ghost-agent.service
[Unit]
Description=Ghost Agent
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ghost-agent
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production
# Allow writing state, compose files, and self-updating the binary
ReadWritePaths=/var/lib/ghost /etc/ghost /usr/local/bin
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
